# CryptoUsersService.service
[Unit]
Description=uWSGI instance to serve UsersService
After=network.target

[Service]
User=CryptoUsersService
Group=www-data
WorkingDirectory=/home/crypto/CryptoUsersService
Environment="PATH=/home/crypto/CryptoUsersService/CryptoUsersServiceEnv/bin"
ExecStart=/usr/local/bin/bin/gunicorn --workers 3 --bind unix:CryptoUsersService.sock -m 007 wsgi:app

[Install]
WantedBy=multi-user.target

